# systemd/eink-rotator.service
[Unit]
Description=E-Ink Dashboard Rotator (Sakura-chan)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Dynamic paths will be replaced by install_service.sh
WorkingDirectory=%WORKING_DIR%
EnvironmentFile=%WORKING_DIR%/.env
Environment=GPIOZERO_PIN_FACTORY=lgpio

# Source the venv and run the dashboard
ExecStart=/bin/bash -lc 'source %WORKING_DIR%/.venv/bin/activate && exec python %WORKING_DIR%/dashboard.py'

# Restart behavior
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eink-rotator

# Security (run as user, not root)
User=%USER%
Group=%USER%

# Access to hardware without being root
SupplementaryGroups=spi gpio

[Install]
WantedBy=default.target