# systemd/eink-rotator.service
[Unit]
Description=E-Ink Dashboard Rotator (Sakura-chan)
After=network-online.target

[Service]
Type=simple
# — Change to your repo location —
WorkingDirectory=/home/tim/Scripting/e-ink_display

# Load environment (for OWM_API_KEY, coords, ROTATE_SECONDS, etc.)
EnvironmentFile=/home/tim/Scripting/e-ink_display/.env

# GPIO backend hint (works with your Waveshare lib)
Environment=GPIOZERO_PIN_FACTORY=lgpio
# Force child dashboards to use the same venv interpreter via env (dashboard.py respects PYTHON_BIN)
Environment=PYTHON_BIN=/home/tim/Scripting/e-ink_display/.venv/bin/python
# Prefer venv bin on PATH for any console scripts
Environment=PATH=/home/tim/Scripting/e-ink_display/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin

# Use your venv’s python
ExecStart=/home/tim/Scripting/e-ink_display/.venv/bin/python /home/tim/Scripting/e-ink_display/dashboard.py

# Make sure the service keeps coming back if it crashes
Restart=always
RestartSec=5



# A little limit bump just in case
NoNewPrivileges=true

[Install]
WantedBy=default.target